# automation/ai.the-debrief.daily-codex-article.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.the-debrief.daily-codex-article</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/papay0/Dev/the-debrief/scripts/run-daily-codex-article.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/papay0/Dev/the-debrief</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ARTICLE_TIMEZONE</key>
    <string>America/Los_Angeles</string>
    <key>HOME</key>
    <string>/Users/papay0</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/papay0/Library/Logs/the-debrief/daily-codex-article.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/papay0/Library/Logs/the-debrief/daily-codex-article.err.log</string>
</dict>
</plist>
